# etc/cron.d/torrentwatch-xa
# /etc/cron.d/torrentwatch-xa: crontab fragment for torrentwatch-xa
0,15,30,45 * * * * www-data /usr/bin/php -q /var/www/html/torrentwatch-xa/twxa_cli.php >/dev/null 2>&1
